# launchd/com.m13v.memory-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.m13v.memory-review</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/matthewdi/user-memories/review/run.sh</string>
	</array>
	<key>StartInterval</key>
	<integer>604800</integer>
	<key>StandardOutPath</key>
	<string>/Users/matthewdi/user-memories/review/logs/launchd-stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/matthewdi/user-memories/review/logs/launchd-stderr.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>HOME</key>
		<string>/Users/matthewdi</string>
	</dict>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
